# systemd/wslptyd.service
# WSL Terminal PTY daemon, managed by systemd (optional — see README.md).
#
# This is a *user* service: it runs as you, so the shells it spawns are yours
# (your uid, your $HOME, your login config). Install it at:
#     ~/.config/systemd/user/wslptyd.service
#
# --persist stops wslptyd from auto-exiting when the last window closes, so the
# service stays resident instead of cycling under Restart=always. KillMode=process
# means a stop/restart only signals the listener, never the cgroup — so it can't
# tear down the connection child and your shells while you're using them.

[Unit]
Description=WSL Terminal PTY daemon (wslptyd) on vsock port 5523
Documentation=https://github.com/robbert1978/WSLterminal
After=default.target
# Don't fall into a failed state from rapid restarts (RestartSec already paces it).
StartLimitIntervalSec=0

[Service]
Type=simple
# Build with native/build.sh and copy artifacts/wslptyd here (NOT /tmp — the
# daemon unlinks /tmp/wslptyd, its staged-copy path, on exit).
ExecStart=%h/.local/bin/wslptyd --vsock 5523 --persist
Restart=always
RestartSec=2
# Only signal the main listener on stop/restart — never the whole cgroup, so
# active sessions (your shells) aren't killed underneath you.
KillMode=process

[Install]
WantedBy=default.target
